# 01-base/archlinux-packages/waveshare-poe-hat-b/waveshare-poe-oled-fan-control.service
[Unit]
Description=waveshare-poe-oled-fan-control

[Service]
ExecStart=/usr/bin/waveshare-poe-oled-fan-control

[Install]
WantedBy=multi-user.target